# PythonScript/WebLaunch/supervisor.conf
[program:web]
command=/root/env/bin/uwsgi --ini uwsgi.ini
directory=/opt
autostart=true
startsecs=0
autorestart=false
redirect_stderr=true
stdout_logfile=/opt/cloudmonitor.log
stopsignal=QUIT